# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.


rm -rf /tmp/extroot.sh
wget https://raw.githubusercontent.com/ittp/openwrt/main/extroot.sh -P /tmp/ 
chmod +x /tmp/extroot.sh
ash /tmp/extroot.sh

exit 0
